# proposals/service_poc/mysrv_proposal2.service
[Unit]
Description=My service 2
After=network.target

[Service]
ExecStart=/opt/srv/srv.py

[Install]
WantedBy=multi-user.target
