# proposals/srv_test/mysrv.service
[Unit]
Description=My service
After=network.target

[Service]
Type=forking
ExecStart=/opt/srv/run.sh

[Install]
WantedBy=multi-user.target
